# _/etc/profile.d/cycript.sh
cyc () { cycript -p $1 /var/root/common.cy > /dev/null ; cycript -p $1; }